# systemd/swing-server.service
[Unit]
# Path to install: /etc/systemd/system/swing-server.service
Description=Swing video upload server
After=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/swing-uploader
ExecStart=/usr/bin/node /home/pi/swing-uploader/server.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
